# Dockers/php-5.0.5/docker.conf
[program:php]
command=php-fpm --fpm-config /etc/php5/php-fpm.conf

[program:memcache]
command=/usr/bin/memcached -u memcache -v

