# Dockers/supervisor.conf
[program:mysql]
command=/root/dockers/start-mysql.sh

[program:php]
command=/root/dockers/start-php.sh

[program:nginx]
command=/root/dockers/start-nginx.sh
